# systemd/gauges.service
[Unit]
Description=Virtual Gauges (PySide6 QML)
After=network.target

[Service]
Type=simple
User=pi
# Set the correct path to the repository directory (example below)
WorkingDirectory=/home/pi/virtual-cluster
ExecStart=/usr/bin/python3 /home/pi/virtual-cluster/src/main.py
Restart=on-failure
Environment=PYTHONUNBUFFERED=1
Environment=TEENSY_DEV=/dev/ttyACM0

[Install]
WantedBy=multi-user.target
